# deploy/lightsail/reins-ci.service
[Unit]
Description=Reins on-box CI/CD (poll main, test in a container, deploy on green)
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
User=ubuntu
# Give a slow npm ci + docker build room to finish.
TimeoutStartSec=1800
ExecStart=/home/ubuntu/reins-src/deploy/lightsail/ci-deploy.sh
